# addons/zbxagentchecker/cron.d/zabbix_agent.cron
# ## Zabbix Agent Checker.
*/10 * * * * root /opt/zabbix/zabbix_agent.sh >/dev/null 2>&1
